# make-snapshots.target
[Unit]
Description=Make Snapshots
Requires=graphical.target
